# ansible_collections/emrys_merlin/uberspace/roles/gitea/templates/gitea.ini.j2
[program:gitea]
directory=%(ENV_HOME)s/var/gitea
command=%(ENV_HOME)s/bin/gitea -c %(ENV_HOME)s/etc/gitea/app.ini web
startsecs=30
autorestart=yes