# ansible_collections/emrys_merlin/uberspace/roles/sr_optimal_edge_bot/templates/sr_optimal_edge_bot.ini.j2
[program:sr_optimal_edge_bot]
environment=DISCORD_TOKEN={{ sr_discord_token }}
command={{ executable }}
autostart=yes
autorestart=yes